# packaging/systemd/krellixd.service.d/bare-metal-hardening.conf
# Drop-in: extra mount-namespace hardening for bare-metal hosts.
#
# These options break /proc visibility on many VPS / container hosts
# (the daemon ends up reporting the underlying host's CPU/memory
# instead of the container's view), so they are NOT enabled by
# default. Only install this drop-in if you are running krellixd on
# real hardware (or a full hypervisor VM with its own kernel and a
# normally-mounted procfs) and want the extra defense in depth.
#
# Install:
#   sudo install -D -m 644 \
#     packaging/systemd/krellixd.service.d/bare-metal-hardening.conf \
#     /etc/systemd/system/krellixd.service.d/bare-metal-hardening.conf
#   sudo systemctl daemon-reload
#   sudo systemctl restart krellixd
#
# Verify krellixd still sees the right CPU/memory afterwards:
#   nc 127.0.0.1 19150 | head -2 | tail -1 \
#     | python3 -c 'import sys,json; o=json.loads(sys.stdin.read()); \
#         print("cpu:", len(o["cpu"]), "mem kB:", o["mem"]["total"])'

[Service]
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
