# packaging/systemd/krellixd.service.d/docker-aliases.conf
# Drop-in: grant the daemon read access to /var/run/docker.sock so it
# can resolve br-<hash> bridge interfaces to their friendly Docker
# network names ("internal", "matrix_internal", etc.). Only needed on
# hosts that actually run Docker — installing this on a host without a
# `docker` group will hard-fail the service start.
#
# Install:
#   sudo install -D -m 644 \
#     packaging/systemd/krellixd.service.d/docker-aliases.conf \
#     /etc/systemd/system/krellixd.service.d/docker-aliases.conf
#   sudo systemctl daemon-reload
#   sudo systemctl restart krellixd

[Service]
SupplementaryGroups=docker
